# contrib/lockdown-daemon.service
[Unit]
Description=Lockdown daemon

[Service]
ExecStart=%h/.local/bin/lockdown daemon
Restart=on-failure

[Install]
WantedBy=default.target
